# docker/netatalk.conf
# set some defaults and start netatalk

[program:netatalk]
command=/usr/sbin/netatalk -d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
